# scripts/news_pipeline/com.100m1s.news-pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.100m1s.news-pipeline</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/seongjinpark/company/100m1s/scripts/news_pipeline/with_lock.sh</string>
        <string>/tmp/100m1s-pipeline.lock</string>
        <string>--</string>
        <string>/bin/bash</string>
        <string>/Users/seongjinpark/company/100m1s/scripts/news_pipeline/pipeline.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/seongjinpark/company/100m1s</string>

    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Minute</key><integer>0</integer></dict>
        <dict><key>Minute</key><integer>10</integer></dict>
        <dict><key>Minute</key><integer>20</integer></dict>
        <dict><key>Minute</key><integer>30</integer></dict>
        <dict><key>Minute</key><integer>40</integer></dict>
        <dict><key>Minute</key><integer>50</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/seongjinpark/company/100m1s/scripts/news_pipeline/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/seongjinpark/company/100m1s/scripts/news_pipeline/logs/launchd-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
