# launchd/drafts/com.100m1s.themes-pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.100m1s.themes-pipeline</string>

    <!--
        REQ-002 패턴 (FLR-2026-04-21): macOS flock 명령 부재 → with_lock.sh(fcntl 래퍼) 경유.
        themes는 메인(/tmp/100m1s-pipeline.lock)과 분리된 별도 락 사용.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/seongjinpark/company/100m1s-pm320/scripts/news_pipeline/with_lock.sh</string>
        <string>/tmp/100m1s-themes.lock</string>
        <string>--</string>
        <string>/bin/bash</string>
        <string>/Users/seongjinpark/company/100m1s-pm320/scripts/themes_pipeline/pipeline.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/seongjinpark/company/100m1s-pm320</string>

    <!--
        메인 파이프라인은 매시 00/10/20/30/40/50 (REQ-001).
        테마뉴스는 +5분 위상 시프트로 writer 충돌 0:
          09:06, 09:16, 09:26, ..., 19:56 (10분 간격)
    -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>6</integer></dict>
        <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>16</integer></dict>
        <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>26</integer></dict>
        <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>36</integer></dict>
        <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>46</integer></dict>
        <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>56</integer></dict>
        <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>6</integer></dict>
        <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>16</integer></dict>
        <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>26</integer></dict>
        <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>36</integer></dict>
        <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>46</integer></dict>
        <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>56</integer></dict>
        <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>6</integer></dict>
        <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>16</integer></dict>
        <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>26</integer></dict>
        <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>36</integer></dict>
        <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>46</integer></dict>
        <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>56</integer></dict>
        <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>6</integer></dict>
        <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>16</integer></dict>
        <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>26</integer></dict>
        <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>36</integer></dict>
        <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>46</integer></dict>
        <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>56</integer></dict>
        <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>6</integer></dict>
        <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>16</integer></dict>
        <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>26</integer></dict>
        <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>36</integer></dict>
        <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>46</integer></dict>
        <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>56</integer></dict>
        <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>6</integer></dict>
        <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>16</integer></dict>
        <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>26</integer></dict>
        <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>36</integer></dict>
        <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>46</integer></dict>
        <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>56</integer></dict>
        <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>6</integer></dict>
        <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>16</integer></dict>
        <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>26</integer></dict>
        <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>36</integer></dict>
        <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>46</integer></dict>
        <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>56</integer></dict>
        <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>6</integer></dict>
        <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>16</integer></dict>
        <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>26</integer></dict>
        <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>36</integer></dict>
        <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>46</integer></dict>
        <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>56</integer></dict>
        <dict><key>Hour</key><integer>17</integer><key>Minute</key><integer>6</integer></dict>
        <dict><key>Hour</key><integer>17</integer><key>Minute</key><integer>16</integer></dict>
        <dict><key>Hour</key><integer>17</integer><key>Minute</key><integer>26</integer></dict>
        <dict><key>Hour</key><integer>17</integer><key>Minute</key><integer>36</integer></dict>
        <dict><key>Hour</key><integer>17</integer><key>Minute</key><integer>46</integer></dict>
        <dict><key>Hour</key><integer>17</integer><key>Minute</key><integer>56</integer></dict>
        <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>6</integer></dict>
        <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>16</integer></dict>
        <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>26</integer></dict>
        <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>36</integer></dict>
        <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>46</integer></dict>
        <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>56</integer></dict>
        <dict><key>Hour</key><integer>19</integer><key>Minute</key><integer>6</integer></dict>
        <dict><key>Hour</key><integer>19</integer><key>Minute</key><integer>16</integer></dict>
        <dict><key>Hour</key><integer>19</integer><key>Minute</key><integer>26</integer></dict>
        <dict><key>Hour</key><integer>19</integer><key>Minute</key><integer>36</integer></dict>
        <dict><key>Hour</key><integer>19</integer><key>Minute</key><integer>46</integer></dict>
        <dict><key>Hour</key><integer>19</integer><key>Minute</key><integer>56</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/seongjinpark/company/100m1s-pm320/scripts/themes_pipeline/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/seongjinpark/company/100m1s-pm320/scripts/themes_pipeline/logs/launchd-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
