# launchd/drafts/com.100m1s.kakao-token-rotate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- 초안 · 미적재 (대표 지시 2026-07-06 "앞으로 스스로 카카오 토큰을 갱신해줘") — 적재는 lead/DevOps 후속:
     cp launchd/com.100m1s.kakao-token-rotate.plist ~/Library/LaunchAgents/
     launchctl bootstrap gui/$UID ~/Library/LaunchAgents/com.100m1s.kakao-token-rotate.plist
     일 1회 08:10 — credential-expiry-watch(08:00) 직후, rotate 결과가 당일 감시·브리핑에 반영되도록 배치 -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.100m1s.kakao-token-rotate</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/seongjinpark/company/100m1s-pm320/scripts/automation/kakao_token_rotate.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/seongjinpark/company/100m1s-pm320</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>10</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/seongjinpark/Library/Logs/100m1s/kakao-token-rotate.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/seongjinpark/Library/Logs/100m1s/kakao-token-rotate.err.log</string>
  <key>RunAtLoad</key>
  <false/>
  <key>KeepAlive</key>
  <false/>
</dict>
</plist>
